# new-api.service
# Production path: /etc/systemd/system/origin-gateway.service
[Unit]
Description=Origin Gateway
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/origin-gateway/data
EnvironmentFile=/opt/origin-gateway/.env
EnvironmentFile=/opt/origin-gateway/native.env
ExecStart=/opt/origin-gateway/current/origin-gateway --port 9320 --log-dir /opt/origin-gateway/logs
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
